# systemd/odoo-backup.service
[Unit]
Description=Odoo backup (DB + filestore, cifrado, local + R2)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/odoo-infrastructure
ExecStart=/usr/bin/docker compose -f backup/docker/docker-compose.yml exec -T backup /usr/local/bin/backup.sh
