# systemd/staging-teardown-boot.service
[Unit]
Description=Teardown incondicional de staging al arrancar el server
Wants=network-online.target
After=docker.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/odoo-infrastructure
ExecStart=/usr/bin/docker compose -f docker/docker-compose.staging.yml down -v

[Install]
WantedBy=multi-user.target
