# systemd/staging-refresh.timer
[Unit]
Description=Corre staging-refresh.service semanalmente

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
